# infra/com.bingo.agent-canvas.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bingo.agent-canvas</string>
  <key>ProgramArguments</key>
  <array>
    <string>NODE_PATH</string>
    <string>SERVER_PATH</string>
  </array>
  <key>WorkingDirectory</key>
  <string>APP_ROOT</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>3</integer>
  <key>StandardOutPath</key>
  <string>LOG_PATH</string>
  <key>StandardErrorPath</key>
  <string>LOG_PATH</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>ENV_PATH</string>
    <key>HOME</key>
    <string>HOME_PATH</string>
  </dict>
</dict>
</plist>
